# faucet.conf
#!upstart
description "faucet.augur.net"

env USER=augur
env GETH_IPC=/home/augur/.ethereum/testnet/geth.ipc

start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start --chdir /home/$USER --chuid $USER --exec /home/$USER/faucet/index.js >> /home/$USER/log/faucet.log 2>&1
