# geth1.conf
#!upstart
description "geth"

env USER=augur

start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start --chdir /home/$USER/.ethereum-augur --chuid $USER --exec /usr/local/bin/geth -- --mine --minerthreads "1" --rpc --rpcapi "eth,net,web3" --ipcapi "admin,db,eth,debug,miner,net,shh,txpool,personal,web3" --networkid 7 --rpccorsdomain "https://dev.augur.net http://dev.augur.net https://augur-dev.firebaseapp.com http://augur-dev.firebaseapp.com http://augur-exp.herokuapp.com https://augur-exp.herokuapp.com https://augur-dev.herokuapp.com http://augur-dev.herokuapp.com http://demo.augur.net http://client.augur.net https://client.augur.net https://augur.firebaseapp.com https://augur-client.herokuapp.com http://localhost:8080 https://localhost:8080 http://127.0.0.1:8080 https://127.0.0.1:8080" --bootnodes "enode://d4f4e7fd3954718562544dbf322c0c84d2c87f154dd66a39ea0787a6f74930c42f5d13ba2cfef481b66a6f002bc3915f94964f67251524696a448ba40d1e2b12@45.33.59.27:30303,enode://a9f34ea3de79cd75ba49c37603d28a7c494f32604b4ad6e3415b4c6020ff5bf38f9772d69362c024355245fe839dd397ff9ec04db70b3258d92259323cb792ae@69.164.196.239:30303,enode://4f23a991ea8739bcc5ab52625407fcfddb03ac31a36141184cf9072ff8bf399954bb94ec47e1f653a0b0fea8d88a67fa3147dbe5c56067f39e0bd5125ae0d1f1@139.162.5.153:30303,enode://bafc7bbaebf6452dcbf9522a2af30f586b38c72c84922616eacad686ab6aaed2b50f808b3f91dba6a546474fe96b5bff97d51c9b062b4a2e8bc9339d9bb8e186@106.184.4.123:30303" --rpcaddr 0.0.0.0 --rpcport 8545 --maxpeers 64 --etherbase 0xaff9cb4dcb19d13b84761c040c91d21dc6c991ec --unlock 0xaff9cb4dcb19d13b84761c040c91d21dc6c991ec --password /home/$USER/.ethereum-augur/.password --datadir /home/$USER/.ethereum-augur >> /home/$USER/log/geth.log 2>&1
