# hauler-geth-main.conf
#!upstart
description "geth"

env USER=jack

start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start --chdir /home/$USER --chuid $USER --exec /usr/local/bin/geth -- --cache 1024 --ipcapi "admin,eth,debug,miner,net,shh,txpool,personal,web3" --maxpeers 1024 --etherbase 0x93d58ab88a403e4c42955a1c1be3e931e12dbe4e --unlock 0x93d58ab88a403e4c42955a1c1be3e931e12dbe4e --password /home/$USER/.ethereum/.password >> /home/$USER/log/geth.log 2>&1
