# morden-parity5.conf
#!upstart
description "parity"

env USER=augur

start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start --chdir /home/$USER --chuid $USER --exec /usr/local/bin/parity -- --chain testnet --jsonrpc --jsonrpc-interface 0.0.0.0 --jsonrpc-cors '*' --cache 2048 --etherbase 0xe434ed7f4684e3d2db25c4937c9e0b7b1faf54c6 >> /home/$USER/log/parity.log 2>&1
