# morden-parity4.conf
#!upstart
description "parity"

env USER=augur

start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start --chdir /home/$USER --chuid $USER --exec /usr/local/bin/parity -- --chain testnet --jsonrpc --jsonrpc-interface 0.0.0.0 --jsonrpc-cors '*' --cache 2048 --etherbase 0x8296eb59079f435275b76058c08b47c4f8965b78 >> /home/$USER/log/parity.log 2>&1
