# peers.conf
#!upstart
description "peers"

env USER=jack

start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start --chdir /home/$USER --chuid $USER --exec /home/$USER/src/evac/peers.js -- -w >> /var/log/peers.log 2>&1
