# slack-irc.conf
#!upstart
description "slack-irc bridge"

env USER=deploy

start on runlevel [2345]
stop on runlevel [016]

respawn

exec start-stop-daemon --start --chuid $USER --exec slack-irc -- --config /home/deploy/slack-irc/config.json 2>> /var/log/slack-irc/slack-irc.log 2>&1
